# luci-theme-cbshield/root/etc/uci-defaults/80_cbshield-theme
#!/bin/sh
# 保持与 Argon 主题兼容，CB-Shield 仅提供业务页面样式

uci -q batch <<'EOF'
set luci.main.mediaurlbase='/luci-static/argon'
set luci.themes.Argon='/luci-static/argon'
commit luci
EOF

exit 0
