# files/etc/uci-defaults/90_v3_optimizations
#!/bin/sh
# CB-Shield first boot defaults (single-network edition)

THEME_PATH='/luci-static/argon'
[ -d /www/luci-static/argon ] || THEME_PATH='/luci-static/bootstrap'

uci -q batch <<EOF
set luci.main.lang='zh_cn'
set luci.main.mediaurlbase='${THEME_PATH}'
set luci.main.homepage='admin/cbshield/dashboard'
commit luci

set system.@system[0].hostname='CB-Shield-Q30'
set system.@system[0].timezone='CST-8'
set system.@system[0].zonename='Asia/Shanghai'
delete system.ntp.server
add_list system.ntp.server='time.cloudflare.com'
add_list system.ntp.server='time.google.com'
add_list system.ntp.server='pool.ntp.org'
commit system

set network.lan.ipaddr='192.168.10.1'
commit network

set dhcp.@dnsmasq[0].cachesize='1500'
commit dhcp
EOF

if [ "$THEME_PATH" = '/luci-static/argon' ]; then
    uci -q set luci.themes.Argon='/luci-static/argon'
else
    uci -q delete luci.themes.Argon
fi
uci commit luci

# Generate wireless config and keep a single LAN-backed dual-band network.
wifi config >/dev/null 2>&1 || wifi detect > /etc/config/wireless
sed -i "s/option country 'US'/option country 'CN'/g" /etc/config/wireless

RADIO0="$(uci -q get wireless.radio0.type)"
RADIO1="$(uci -q get wireless.radio1.type)"

if [ -n "$RADIO0" ]; then
    uci -q batch <<'EOF'
set wireless.radio0.channel='auto'
set wireless.radio0.band='2g'
set wireless.radio0.htmode='HE40'
set wireless.radio0.country='CN'
set wireless.radio0.disabled='0'
set wireless.default_radio0.disabled='0'

set wireless.office_24g=wifi-iface
set wireless.office_24g.device='radio0'
set wireless.office_24g.mode='ap'
set wireless.office_24g.ssid='CB-Shield-Office'
set wireless.office_24g.network='lan'
set wireless.office_24g.encryption='psk2'
set wireless.office_24g.key='CBShield@Office2024'
EOF
fi

if [ -n "$RADIO1" ]; then
    uci -q batch <<'EOF'
set wireless.radio1.channel='auto'
set wireless.radio1.band='5g'
set wireless.radio1.htmode='HE80'
set wireless.radio1.country='CN'
set wireless.radio1.disabled='0'
set wireless.default_radio1.disabled='0'

set wireless.office_5g=wifi-iface
set wireless.office_5g.device='radio1'
set wireless.office_5g.mode='ap'
set wireless.office_5g.ssid='CB-Shield-Office-5G'
set wireless.office_5g.network='lan'
set wireless.office_5g.encryption='psk2'
set wireless.office_5g.key='CBShield@Office2024'
EOF
fi

uci commit wireless
wifi reload >/dev/null 2>&1 || true

# Mark wizard required if root password is missing.
mkdir -p /etc/cbshield
ROOT_HASH="$(grep '^root:' /etc/shadow | cut -d: -f2)"
if [ -z "$ROOT_HASH" ] || [ "$ROOT_HASH" = "*" ] || [ "$ROOT_HASH" = "!" ]; then
    touch /etc/cbshield/wizard_required
    uci -q set cb-wizard.main.required='1'
else
    rm -f /etc/cbshield/wizard_required
    uci -q set cb-wizard.main.required='0'
fi
uci commit cb-wizard

# Enable CB ops services.
for svc in cb-riskcontrol cb-portal cb-healthcheck cb-ha-monitor cb-dns-guard; do
    /etc/init.d/$svc enable 2>/dev/null
    /etc/init.d/$svc start 2>/dev/null
done

/usr/bin/cb-eventlog write system "single-network defaults applied"

exit 0
