# files/etc/rc.local
#!/bin/sh
# ============================================================
# CB-Shield-Q30 启动脚本
# ============================================================

# 加载 BBR 模块
modprobe tcp_bbr 2>/dev/null

# 应用 sysctl 优化参数
sysctl -p /etc/sysctl.conf 2>/dev/null

# 应用策略路由
ip rule add fwmark 0x30 table 30 2>/dev/null
ip route add default via $(uci get network.wan.gateway 2>/dev/null || echo "0.0.0.0") table 30 2>/dev/null

# 启动风控服务
/etc/init.d/cb-riskcontrol start 2>/dev/null

exit 0
